# provisioning/etc/systemd/system/followback.timer
[Unit]
Description=run followback every day

[Timer]
OnUnitActiveSec=1d

[Install]
WantedBy=multi-user.target
